# scripts/sai-jig.service
[Unit]
Description=Sai JIG

[Service]
ExecStart=/usr/local/bin/sai-jig

[Install]
WantedBy=multi-user.target

